# pulling_from_lux_sensor/crontab
*/1 * * * * /usr/local/bin/python /app/lux_data_extraction_job.py > /proc/1/fd/1 2>/proc/1/fd/2
